# scripts/supervisord.conf
; Tent OS Supervisor 配置 —— all-in-one 模式
; 在一个容器内运行所有组件

[supervisord]
nodaemon=true
user=root
logfile=/app/logs/supervisord.log
pidfile=/tmp/supervisord.pid

[program:memory-worker]
command=python -m tent_os.cli worker memory --config %(ENV_TENT_OS_CONFIG)s
autostart=true
autorestart=true
stdout_logfile=/app/logs/memory-worker.log
stderr_logfile=/app/logs/memory-worker.log
priority=10

[program:governance-worker]
command=python -m tent_os.cli worker governance --config %(ENV_TENT_OS_CONFIG)s
autostart=true
autorestart=true
stdout_logfile=/app/logs/governance-worker.log
stderr_logfile=/app/logs/governance-worker.log
priority=20
startsecs=3

[program:scheduler-worker]
command=python -m tent_os.cli worker scheduler --config %(ENV_TENT_OS_CONFIG)s
autostart=true
autorestart=true
stdout_logfile=/app/logs/scheduler-worker.log
stderr_logfile=/app/logs/scheduler-worker.log
priority=20
startsecs=3

[program:api-server]
command=python -m tent_os.cli server --host 0.0.0.0 --port 8000 --config %(ENV_TENT_OS_CONFIG)s
autostart=true
autorestart=true
stdout_logfile=/app/logs/api-server.log
stderr_logfile=/app/logs/api-server.log
priority=30
startsecs=5
